# rootfs/etc/s6-overlay/s6-rc.d/svc-stdout/run
#!/usr/bin/with-contenv bash

set -e

dockerize -stdout /var/log/nginx/error.log -poll
